# unison
#!/sbin/openrc-run

# Copyright (c) 2026 Wizarsy
# 
# Permission is hereby granted, free of charge, to any person obtaining a copy
# of this software and associated documentation files (the "Software"), to deal
# in the Software without restriction, including without limitation the rights
# to use, copy, modify, merge, publish, distribute, sublicense, and/or sell
# copies of the Software, and to permit persons to whom the Software is
# furnished to do so, subject to the following conditions:
# 
# The above copyright notice and this permission notice shall be included in all
# copies or substantial portions of the Software.
# 
# THE SOFTWARE IS PROVIDED "AS IS", WITHOUT WARRANTY OF ANY KIND, EXPRESS OR
# IMPLIED, INCLUDING BUT NOT LIMITED TO THE WARRANTIES OF MERCHANTABILITY,
# FITNESS FOR A PARTICULAR PURPOSE AND NONINFRINGEMENT. IN NO EVENT SHALL THE
# AUTHORS OR COPYRIGHT HOLDERS BE LIABLE FOR ANY CLAIM, DAMAGES OR OTHER
# LIABILITY, WHETHER IN AN ACTION OF CONTRACT, TORT OR OTHERWISE, ARISING FROM,
# OUT OF OR IN CONNECTION WITH THE SOFTWARE OR THE USE OR OTHER DEALINGS IN THE
# SOFTWARE.

supervisor="supervise-daemon"
description="Starts unison profile files"
profile="${RC_SVCNAME#*.}"

command="/usr/bin/unison"
command_args="${profile} -batch -killserver -log -logfile ${profile}.log -repeat watch+${watch_scan_interval:-1800}"

depend() {
  [ "${rc_need+set}" = "set" ] && return 0
  if "$(grep -qE "^root ?= ?(ssh|socket 2> /dev/null)" "$profile")"; then
    need net
    ewarn "You have configured a remote root in ${profile}.prf"
    ewarn "You need to add rc_need=\"net.FOO\" to your ${HOME}/.config/rc/conf.d/${RC_SVCNAME}"
    ewarn "where FOO is your network interface"
  fi
}

start_pre() {
  if [ "$profile" = "$RC_SVCNAME" ]; then
    eerror "${RC_SVCNAME} cannot be started directly. You must create" 
    eerror "symbolic links to it for the profiles you want to start"
    return 1
  elif [ ! -e "${HOME}/.unison/${profile}" ]; then
    eerror "Profile ${profile}.prf not found in ${HOME}/.unison"
    return 1
  fi
}